# markermate.service
[Unit]
Description=MarkerMate Backend API
After=network.target

[Service]
Type=simple
User=ryan
WorkingDirectory=/home/ryan/marker-mate/backend
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
